# systemd/.config/systemd/user/hyprland.service
[Unit]
Description=%p
BindsTo=graphical-session.target

[Service]
Type=notify
EnvironmentFile=%h/.env
ExecStart=/usr/bin/Hyprland
