# systemd/.config/systemd/user/xdg-auto-apps.service
[Unit]
Description="Imperative shit"
Wants=xdg-desktop-autostart.target

[Service]
ExecStart=sh -c "while [ 1 ];do sleep 1;done"
