# systemd/.config/systemd/user/xdg-auto-apps.timer
[Unit]
Description=Run autostart desktop apps

[Timer]
OnActiveSec=5s
